# sqls/move2docker/etc/rc.local
#!/bin/bash
chmod 222 /dev/pts/0
crond
service mysql start
curl -d "restart backend" http://172.33.1.1:26410
sleep infinity
