# ndbd/move2docker/etc/rc.local
#!/bin/bash
set -x
if [ -d "/etc/notfirstrun" ] 
then
    exec ndbd --skip-daemon
else
    mkdir -p /etc/notfirstrun
    exec ndbd --skip-daemon --initial
fi

